# samples/plakar-hiking.timer
[Unit]
Description=Backup hiking
[Timer]
OnCalendar=*-*-* 23:40:00
[Install]
WantedBy=timers.target
