# samples/plakar-hiking.service
[Unit]
Description=Backup hiking
[Service]
User=roelof
ExecStart=/usr/bin/plakar at @roelof backup -ignore-file /data/plakar/exclude.txt @hiking
